# deploy/mdd-candy.service
# systemd unit for the MDD Candy Node server.
# Install:
#   sudo cp deploy/mdd-candy.service /etc/systemd/system/mdd-candy.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now mdd-candy
# Logs:
#   journalctl -u mdd-candy -f
#
# Adjust User, WorkingDirectory and the npm path to match your VPS.

[Unit]
Description=MDD Candy (Node API + static frontend)
After=network.target

[Service]
Type=simple
User=timmi
WorkingDirectory=/home/timmi/projects/MDD_Candy/server
EnvironmentFile=/home/timmi/projects/MDD_Candy/server/.env
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
# Node respects PORT/HOST from the EnvironmentFile above.

[Install]
WantedBy=multi-user.target
